# com.audiotools.normalizer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.audiotools.normalizer</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/simondavis/projects/audio-normalisation/watcher.py</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardErrorPath</key>
    <string>/Users/simondavis/Library/Logs/audio-normalizer-error.log</string>

    <key>StandardOutPath</key>
    <string>/Users/simondavis/Library/Logs/audio-normalizer-output.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/simondavis/projects/audio-normalisation</string>
</dict>
</plist>
